# roles/teeworlds/templates/teeworlds.service.j2
[Unit]
Description=teeworlds server

[Service]
ExecStart="{{ teeworlds__install_dir }}/current/teeworlds_srv -f server.cfg"

[Install]
WantedBy=multi-user.target